# fsd-toggle.service
[Unit]
Description=Tesla FSD Toggle Server (Comma 3 / HW3)
# Start after /data is mounted (scripts and venv paths live there).
After=network.target local-fs.target
RequiresMountsFor=/data

[Service]
Type=simple
# Match openpilot: panda lives under the openpilot tree; venv has usb1/opendbc.
Environment=PYTHONPATH=/data/pythonpath
ExecStart=/usr/local/venv/bin/python3 /data/fsd_toggle_server.py
WorkingDirectory=/data
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
