# fsd-toggle.service
[Unit]
Description=Tesla FSD Toggle Server (Comma 3 / HW3)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /data/fsd_toggle_server.py
WorkingDirectory=/data
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
